# supervisor/conf.d/django.conf
[program:django-app]
directory=/usr/local/kube_cd_project
command=DJANGO_COMMAND
