# supervisor/conf.d/django.conf
[program:django-app]
directory=/usr/local/kube_cd_project
command=gunicorn -b 0.0.0.0:8000 kube_cd_project.wsgi
